# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/chargelink/supervisord.log
logfile_maxbytes=10MB
pidfile=/run/supervisord.pid

; ── Spring Boot backend ────────────────────────────────────────────────────────
[program:backend]
command=java -jar /app/backend.jar
directory=/app
user=chargelink
autostart=true
autorestart=true
startsecs=15
startretries=3
stdout_logfile=/var/log/chargelink/backend-console.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/chargelink/backend-console.log
stderr_logfile_maxbytes=50MB

; ── nginx frontend ─────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/chargelink/nginx-access.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/chargelink/nginx-error.log
stderr_logfile_maxbytes=10MB
